# .config/autostart/FB_messenger_chrome.desktop
#!/usr/bin/env xdg-open
[Desktop Entry]
Version=1.0
Terminal=false
Type=Application
Name=Messenger
Exec=/opt/google/chrome/google-chrome --profile-directory=Default --app-id=bbdeiblfgdokhlblpgeaokenkfknecgl --app-run-on-os-login-mode=windowed
Icon=chrome-bbdeiblfgdokhlblpgeaokenkfknecgl-Default
StartupWMClass=crx_bbdeiblfgdokhlblpgeaokenkfknecgl
X-GNOME-Autostart-Delay=7
