# .config/autostart/messages.desktop
#!/usr/bin/env xdg-open
[Desktop Entry]
Version=1.0
Terminal=false
Type=Application
Name=Messages
Exec=/opt/google/chrome/google-chrome --profile-directory=Default --app-id=hpfldicfbfomlpcikngkocigghgafkph
Icon=chrome-hpfldicfbfomlpcikngkocigghgafkph-Default
StartupWMClass=crx_hpfldicfbfomlpcikngkocigghgafkph
X-GNOME-Autostart-Delay=9
